# scripts/systemd/chrome-cdp.service
[Unit]
Description=Chrome with Remote Debugging (CDP)
After=xfce-desktop.service

[Service]
Type=simple
Environment=DISPLAY=:1
ExecStart=/usr/bin/google-chrome \
    --remote-debugging-port=9222 \
    --no-first-run \
    --no-default-browser-check \
    --disable-default-apps \
    --disable-gpu \
    --window-size=1280,720 \
    --user-data-dir=%h/.chrome-cdp-profile \
    about:blank
ExecStop=/bin/kill -TERM -$MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
